# autonomous-workbench/infra/systemd/openclaw-gateway.service.d/alert.conf
[Unit]
OnFailure=openclaw-gateway-alert.service
